# notes-editor.service
[Unit]
Description=Daily Notes Editor (Go Server)
After=network.target

[Service]
Type=simple
User=dev
WorkingDirectory=/home/dev/notes_editor/server
EnvironmentFile=/home/dev/notes_editor/server/.env
Environment="PATH=/usr/local/go/bin:/home/dev/.local/bin:/usr/local/bin:/usr/bin:/bin"
StandardOutput=journal
StandardError=journal
ExecStart=/home/dev/notes_editor/server/bin/server
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
