# notes-editor.service
[Unit]
Description=Daily Notes Editor
After=network.target

[Service]
Type=simple
User=dev
WorkingDirectory=/home/dev/notes_editor
Environment="PATH=/home/dev/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
StandardOutput=journal
StandardError=journal
ExecStart=/home/dev/notes_editor/.venv/bin/uvicorn server.web_app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
