# notes-pi-gateway.service
[Unit]
Description=Notes Editor Gateway Sidecar
After=network.target
Wants=network.target

[Service]
Type=simple
User=dev
Group=dev
WorkingDirectory=/home/dev/notes_editor/pi-gateway
EnvironmentFile=/home/dev/notes_editor/server/.env
Environment=PI_GATEWAY_PORT=4317
Environment=PI_GATEWAY_MODE=pi_rpc
Environment=PI_GATEWAY_PI_PROVIDER=anthropic
Environment=NOTES_SERVER_URL=http://127.0.0.1:80
ExecStart=/home/dev/.local/node20/bin/node /home/dev/notes_editor/pi-gateway/dist/server.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
