# notes-pi-gateway.service
[Unit]
Description=Notes Editor Gateway Sidecar
After=network.target
Wants=network.target

[Service]
Type=simple
User=dev
Group=dev
WorkingDirectory=/home/dev/notes_editor/pi-gateway
Environment=PI_GATEWAY_PORT=4317
Environment=PI_GATEWAY_MODE=claude_cli
Environment=PI_GATEWAY_CLAUDE_BIN=/home/dev/.local/bin/claude
ExecStart=/usr/bin/node /home/dev/notes_editor/pi-gateway/dist/server.js
Restart=always
RestartSec=3
StandardOutput=append:/home/dev/notes_editor/pi-gateway/gateway.log
StandardError=append:/home/dev/notes_editor/pi-gateway/gateway.err.log

[Install]
WantedBy=multi-user.target
